# rootfs/s6/spigot-list-controller/run
#!/bin/sh
exec spigot-exec spigot-list-controller
