# rootfs/s6/spigot-metrics-controller/run
#!/bin/sh
exec spigot-exec spigot-metrics-controller
