# rootfs/s6/spigot-stop-controller/run
#!/bin/sh
exec spigot-exec spigot-stop-controller
