# rootfs/etc/s6/spigot-log/run
#!/bin/sh
exec s6-setuidgid craftbukkit /spigot-log.sh
