# rootfs/s6/spigot-backup-controller/run
#!/bin/sh
exec spigot-exec spigot-backup-controller
