# rootfs/s6/darkhttpd/run
#!/bin/bash

set -euo pipefail

source /usr/local/etc/spigot.conf

exec darkhttpd \
"${METRICS_DIR}" \
"${METRICS_DARKHTTPD_FLAGS[@]}"
